# systemd/security-retention.service
[Unit]
Description=Security Retention Cleanup Run
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/security-scanning-platform
ExecStart=/opt/security-scanning-platform/scripts/schedule_retention.sh
User=securityscan
Group=securityscan
