# systemd/security-scanner.service
[Unit]
Description=Security Scanner Batch Run
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/security-scanning-platform
ExecStart=/opt/security-scanning-platform/scripts/schedule_scan.sh
User=securityscan
Group=securityscan

[Install]
WantedBy=multi-user.target
